# linux/arch-workstation/dotfiles/bash/profile.d/smoothoperator-workstation.sh
# SmoothOperator workstation — ensure user-local tools on PATH
case ":${PATH}:" in
  *":${HOME}/.local/bin:"*) ;;
  *) export PATH="${HOME}/.local/bin:${PATH}" ;;
esac
